# local-node/setup/local-node.service
[Unit]
Description=Local Node
After=network.target

[Service]
ExecStart=python /home/pi/code/local-node.py
WorkingDirectory=/home/pi/code
StandardOutput=inherit
StandardError=inherit
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
